# deployment/supervisor-discordbot.conf
[program:fpsociety-discordbot]
process_name=%(program_name)s
command=php /var/www/html/artisan discordbot:start
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/discord-bot.log
stopwaitsecs=10
stopsignal=TERM
startsecs=3
